# etc/systemd/system/on-calendar.timer
[Unit]
Description=Test for OnCalendar timer

[Timer]
OnCalendar=*-*-* *:00:00
AccuracySec=1s

[Install]
WantedBy=timers.target
