# etc/systemd/system/pcchecker_alarm.timer
[Unit]
Description=PcChecker timer

[Timer]
OnCalendar=*-*-* *:*:*
AccuracySec=1s

[Install]
WantedBy=timers.target
