# etc/systemd/system/pcchecker_alarm.timer
[Unit]
Description=PcChecker timer

[Timer]
OnCalendar=*-*-* *:00:00
AccuracySec=1s

[Install]
WantedBy=timers.target
